# switch/switch.crontab
* * * * *  /home/r00ts/util/update_switch_status.py
