# switch/switch.crontab
* * * * *  python /home/r00ts/update_switch_status.py --home_id `cat /home/r00ts/switchconfigs/r00tzRegistered | tr -d '"'` --get --switch_id `cat /home/r00ts/swith/configs/r00tzSwitchID | tr -d '"'` --updatestatus
*/15 * * * *  python /home/r00ts/update_switch_status.py  --update --apply
